# contrib/systemd/solus-server.service
[Unit]
Description=Solus AI Workflow Engine — Web Server
After=network.target ollama.service
Wants=solus-worker.service

[Service]
Type=simple
User=solus
Group=solus
ExecStart=/usr/bin/env solus serve --host 127.0.0.1 --port 8765
Restart=on-failure
RestartSec=5

# Paths
Environment=HOME=/home/solus
WorkingDirectory=/home/solus

# Hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/solus/.local/share/solus /home/solus/.config/solus
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
